# backend/deploy/agricultural-api.service
[Unit]
Description=Agricultural Data Platform API
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/agricultural-api/backend
Environment="PATH=/var/www/agricultural-api/backend/venv/bin"
EnvironmentFile=/var/www/agricultural-api/.env
ExecStart=/var/www/agricultural-api/backend/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 4
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target